# docker/rskj/supervisord.conf
[supervisord]
nodaemon=true

[program:rskj]
command=/usr/bin/java -Drpc.providers.web.ws.port=8546 -Drpc.providers.web.http.port=8545 -Dblockchain.config.hardforkActivationHeights.iris300=1 -Dminer.client.autoMine=true -Drpc.providers.web.ws.enabled=true -cp /usr/share/rsk/rsk.jar co.rsk.Start > /dev/null 2>&1 &
autostart=true
user=rsk